# playbooks/files/autoborg.timer
[Unit]
Description=Run borgmatic backup

[Timer]
OnCalendar=*-*-* 23:15:00
Persistent=true

[Install]
WantedBy=timers.target
