# playbooks/files/backup_snapshot.timer
[Unit]
Description=Run backup_snapshots.sh daily at 23:30:00

[Timer]
OnCalendar=*-*-* 23:30:00
Persistent=true

[Install]
WantedBy=timers.target
